# comon/rootfs/etc/services.d/comon/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: CoMon
# Starts CoMon
# ==============================================================================

## Run CoMon
cd /opt/comon || bashio::exit.nok "Could not change directory to CoMon"

bashio::log.info "Starting CoMon..."

exec .//CoMon.Web.Host